# systemd/docker-volume-helpers-upgrade.service
[Unit]
Description=Upgrade Docker volume helpers

[Service]
Type=oneshot
TimeoutSec=5m
ExecStart=/usr/local/bin/docker-volume-helpers-upgrade
